# client/terminaliso/iso/airootfs/etc/systemd/system/getty@tty3.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin irssi --noclear %I 38400 linux
