# client/terminaliso/iso/airootfs/etc/systemd/system/sasprepare.service
[Unit]
Description=SaS Preparation Script

[Service]
ExecStart=/root/sasprepare.sh

[Install]
WantedBy=multi-user.target
